# supervisord.ini
[inet_http_server]
port=127.0.0.1:9001        ; (ip_address:port specifier, *:port for all iface)

[supervisord]
logfile=/tmp/supervisord.log ; (main log file;default $CWD/supervisord.log)
logfile_maxbytes=50MB        ; (max main logfile bytes b4 rotation;default 50MB)
logfile_backups=10           ; (num of main logfile rotation backups;default 10)
loglevel=info                ; (log level;default info; others: debug,warn,trace
)
pidfile=/service/supervisord.pid ; (supervisord pidfile;default supervisord.pid)
nodaemon=true               ; (start in foreground if true;default false)
minfds=1024                  ; (min. avail startup file descriptors;default 1024
)
minprocs=200                 ; (min. avail process descriptors;default 200)

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///tmp/supervisor.sock ; use a unix:// URL  for a unix socket

[program:notify]
;command = python /Users/guoyucheng/master/geser/notify/manage.py service_runner
command = python /service/manage.py service_runner
process_name=%(program_name)s_%(process_num)02d
user=root
numprocs=4                   ; å¯å¨å ä¸ªè¿ç¨
autostart=true                ; éçsupervisordçå¯å¨èå¯å¨
autorestart=true              ; èªå¨éå¯ããå½ç¶è¦éä¸äº
#startretries=10               ; å¯å¨å¤±è´¥æ¶çæå¤éè¯æ¬¡æ°
exitcodes=0                 ; æ­£å¸¸éåºä»£ç ï¼æ¯è¯´éåºä»£ç æ¯è¿ä¸ªæ
stopsignal=KILL               ; ç¨æ¥ææ­»è¿ç¨çä¿¡å·
stopwaitsecs=10               ; åéSIGKILLåçç­å¾æ¶é´
;stdout_logfile=/Users/guoyucheng/master/geser/notify/out.log
;stdout_logfile=/service/notify/out.log
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
